# weather_pull.service
[Unit]
Description=Darksky API fetch

[Service]
User=%I
Restart=always
RestartSec=1800s
ExecStart=/usr/bin/python3 /home/pi/weather_station/pull.py
